# deploy/avril-rag.service
[Unit]
Description=Avril RAG - Automated News Pipeline for Agentic Business Ideas
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=backend
Group=backend
WorkingDirectory=/srv/avril
ExecStart=/usr/bin/python3 /srv/avril/main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=avril-rag
SyslogFacility=local7

# Seguridad
NoNewPrivileges=true
ProtectHome=yes
ProtectSystem=strict

# Timeouts (30 min para el pipeline)
TimeoutStartSec=30m
Restart=no

[Install]
WantedBy=multi-user.target


